# systemd/dht-sensor.service
[Unit]
Description=Smart frame — DHT22 room sensor reader
After=multi-user.target

[Service]
ExecStart=/usr/bin/python3 /home/pi/raspi-4b/scripts/dht_read.py
Restart=always
User=pi

[Install]
WantedBy=multi-user.target
